# roles/envoy-install/templates/envoy-hot-restart.service.j2
# {{ ansible_managed }}

[Unit]
Description=Envoy Proxy with Hot Restart
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ envoy_user }}
Group={{ envoy_group }}
ExecStart=/usr/local/bin/hot-restarter.py /usr/local/bin/start-envoy.sh
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hot restart signals
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/envoy /var/log/envoy_admin.log
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target 
